# systemd/vr-stack-control.service
[Unit]
Description=VR Stack Control (tracking -> server -> VR app)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Start the stack in the correct order
ExecStart=%h/.local/bin/vr-stack-run.sh start

# Stop cleanly using the same config/patterns
ExecStop=%h/.local/bin/vr-stack-run.sh stop

TimeoutStopSec=15

[Install]
WantedBy=default.target
